# scripts/active/com.unaligned.dailyscraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.unaligned.dailyscraper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/asherweisberger/.local/bin/run_codex_daily_scraper.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/asherweisberger</string>
  <key>StandardOutPath</key>
  <string>/Users/asherweisberger/Library/Logs/unaligned/dailyscraper.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/asherweisberger/Library/Logs/unaligned/dailyscraper.launchd.error.log</string>
</dict>
</plist>